# vps/etc/supervisor/conf.d/oauth2_proxy.conf
[program:oauth2_proxy]
command=/usr/local/oauth2_proxy/oauth2_proxy -config=/usr/local/oauth2_proxy/oauth2_proxy.cfg
user=www-data
authstart=true
autorestart=true
directory=/usr/local/oauth2_proxy