# vps.old/etc/supervisor/conf.d/caddy.conf
[program:caddy]
command=/opt/caddy/caddy -conf="/etc/caddy/Caddyfile" -agree -log "/var/log/caddy/caddy.log"
environment=HOME="/etc/caddy"
user=www-data
autostart=true
autorestart=true
