# vps.old/etc/supervisor/conf.d/gateway-client.conf
[program:gateway-client]
command=/usr/local/gateway-client/gateway-client
user=www-data
autostart=true
autorestart=true
directory=/usr/local/gateway-client